# infra/systemd/openclaw-brain.service
[Unit]
Description=OpenClaw Brain knowledge updater heartbeat
Documentation=file:/opt/openclaw-brain-v2/BOTS.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/openclaw-brain-v2
EnvironmentFile=-/opt/openclaw-brain-v2/.env
Environment=OPENCLAW_ROOT=/opt/openclaw-brain-v2
Environment=BOT_HEARTBEAT_INTERVAL_SECONDS=300
ExecStart=/usr/bin/python3 /opt/openclaw-brain-v2/knowledge_updater/schedulers/heartbeat_loop.py
Restart=always
RestartSec=15
User=openclaw
Group=openclaw
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/openclaw-brain-v2/knowledge /opt/openclaw-brain-v2/knowledge_updater /opt/openclaw-brain-v2/reports /opt/openclaw-brain-v2/.openclaw

[Install]
WantedBy=multi-user.target
